# roles/caddy_server/templates/caddy-updater.service.j2
[Unit]
Description=Caddy updater service

[Service]
Type=oneshot
ExecStart={{ caddy_upgrade_script_path }}
